# systemd/claude-bot.service
[Unit]
Description=Claude Instance Manager Telegram Bot
After=network.target

[Service]
Type=simple
# TODO: Replace YOUR_USERNAME with the account that will run the bot
User=YOUR_USERNAME
# TODO: Replace /opt/claude-instance-manager with your actual install path
WorkingDirectory=/opt/claude-instance-manager
ExecStart=/opt/claude-instance-manager/.venv/bin/python -m src.main /opt/claude-instance-manager/config.yaml
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
